# systemd/aegis.service
[Unit]
Description=Aegis OTA daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/aegis daemon --config /etc/aegis/system.conf
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
